# devops/chef/cookbooks/cookbook-tomcat/templates/default/tomcat.service.erb
# Systemd unit file for tomcat
[Unit]
Description=Apache Tomcat Web Application Container
After=syslog.target network.target tomcat-keygen.service
Wants=tomcat-keygen.service

[Service]
Type=forking

EnvironmentFile=-/etc/sysconfig/tomcat
ExecStart=/opt/tomcat/bin/startup.sh
ExecStop=/bin/kill -15 $MAINPID

User=tomcat
Group=tomcat

[Install]
WantedBy=multi-user.target
